# adbyss/skel/systemd/adbyss.service
[Unit]
Description=Adbyss

[Service]
Type=oneshot
ExecStart=/usr/bin/adbyss --systemd

[Install]
WantedBy=multi-user.target
